# server/rewind.service
[Unit]
Description=Rewind sailing game server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/home/ubuntu/rewind http
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/ubuntu/rewind.env

[Install]
WantedBy=multi-user.target
